# .config/autostart/ibus.desktop
[Desktop Entry]
Exec=ibus-daemon -drx
Icon=system-run
StartupNotify=true
Terminal=false
Type=Application
